# emqx/rootfs/etc/s6-overlay/s6-rc.d/init-emqx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: EMQX
# Configures EMQX
# ==============================================================================
mkdir -p \
  /data/emqx/data \
  /data/emqx/etc \
  /data/emqx/plugins
